# voz/launchagents/com.satori.voz.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.satori.voz.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/lucianopablolp/Documents/Claude/Projects/SatoriOS/voz/agent/.venv/bin/python</string>
    <string>serve_voz.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/lucianopablolp/Documents/Claude/Projects/SatoriOS/voz/web</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>
  <key>StandardOutPath</key>
  <string>/Users/lucianopablolp/Library/Logs/satori-voz-server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lucianopablolp/Library/Logs/satori-voz-server.log</string>
</dict>
</plist>
